# etc/ft.conf
[program:ft]
command = /home/jpm/from-traccar/ft-run.sh
directory = /tmp/

startretries = 20
autorestart = true

user = jpm
environment= LOGNAME="jpm", USER="jpm"

stdout_logfile = /tmp/ft.log
